# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command = /bin/bash -c '/entrypoint.sh && exec /usr/sbin/nginx -g "daemon off;"'

[program:crond]
command = /usr/sbin/crond -f
autostart=true

