# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c '/entrypoint.sh && /usr/sbin/nginx -g "daemon off;"'

[program:crond]
command = /usr/sbin/cron
user = root
autostart = true

